# lib/farmstead/scaffold/supervisord.conf.erb
[supervisord]
nodaemon=true

[program:producer]
command=/usr/local/bin/ruby /service/project.rb %(ENV_SERVICE)s producer
redirect_stderr=true

[program:consumer]
command=/usr/local/bin/ruby /service/project.rb %(ENV_SERVICE)s consumer
redirect_stderr=true
