# lib/farmstead/scaffold/supervisord.conf.erb
[supervisord]
nodaemon=true

[program:service]
command=/usr/local/bin/ruby /service/project.rb %(ENV_SERVICE)s Service
redirect_stderr=true

[program:api]
command=/usr/local/bin/ruby /service/project.rb %(ENV_SERVICE)s API
redirect_stderr=true